# setup/supervisor.mud.conf
[program:mud]
command = gunicorn --config app/gunicorn_settings.py app:app
directory = /vagrant
user = vagrant
